# conftest/testcrash.conf
# test program: crashes every 10 seconds
[program:testcrash]
command=/R/bin/crashin 10

autostart=true
autorestart=true

stdout_logfile=/L/testcrash.log
redirect_stderr=true
